# treasure-map-bot.service
[Unit]
Description=Treasure Map Telegram Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/workspace
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/python3 /workspace/bot-robust.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=treasure-map-bot

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/workspace/uploads /workspace/treasure_map.db /workspace/bot.log

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target